# ansible/roles/result/files/result.service
[Unit]
Description=Result Node.js App

[Service]
Environment=POSTGRES_USER={{ postgresql_username }}
Environment=POSTGRES_PASSWORD={{ postgresql_password }}
Environment=POSTGRES_HOST={{ hostvars['postgresql'].ansible_host }}
Environment=POSTGRES_PORT={{ postgresql_port }}
Environment=POSTGRES_DB={{ postgresql_db }}
WorkingDirectory=/var/www/result
ExecStart=node server.js
Restart=on-failure
